# flow-tracker/scripts/plists/com.flowtracker.daily-fno.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily F&O bhavcopy + participant OI fetch. Runs at 18:30 IST (after NSE
  publishes end-of-day F&O files around 18:00). Weekends run too — they
  will no-op when NSE has no new data (CLI is idempotent).

  This is a template: __HOME__ is substituted with $HOME by
  scripts/install-launch-agents.sh at install time. Do not copy this file
  to ~/Library/LaunchAgents/ directly — run the install script instead.

  Install:
      bash scripts/install-launch-agents.sh

  Verify:
      launchctl list | grep flowtracker.daily-fno

  Uninstall:
      launchctl unload ~/Library/LaunchAgents/com.flowtracker.daily-fno.plist
      rm ~/Library/LaunchAgents/com.flowtracker.daily-fno.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flowtracker.daily-fno</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.local/share/flowtracker/scripts/daily-fno.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/flowtracker/cron.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/flowtracker/cron.log</string>
</dict>
</plist>
